# zabbix/rsyslogd.conf
[supervisord]
nodaemon=true

[program:rsyslog]
priority=1005
command = bash -c "/usr/sbin/rsyslogd -i /var/run/syslogd.pid"
startsecs=100
autorestart=true
startretries=10
